# rules/persistence_boot_or_logon_autostart_execution.yml
- group: Registry Run Keys / Startup Folder
  description: |
    Adversaries may achieve persistence by adding a program to a startup
    folder or referencing it with a Registry run key. Adding an entry to
    the run keys in the Registry or startup folder will cause the program
    referenced to be executed when a user logs in. These programs will be
    executed under the context of the user and will have the account's
    associated permissions level.
  labels:
    tactic.id: TA0006
    tactic.name: Persistence
    tactic.ref: https://attack.mitre.org/tactics/TA0006/
    technique.id: T1547
    technique.name: Boot or Logon Autostart Execution
    technique.ref: https://attack.mitre.org/techniques/T1547/
    subtechnique.id: T1547.001
    subtechnique.name: Registry Run Keys / Startup Folder
    subtechnique.ref: https://attack.mitre.org/techniques/T1547/001/
  rules:
    - name: Unusual file written or modified in Startup folder
      description: |
        Identifies suspicious files written to the startup folder that would
        allow adversaries to maintain persistence on the endpoint.
      condition: >
        create_file
            and
          (
            file.extension in ('.vbs', '.js', '.jar', '.exe', '.dll', '.com', '.ps1', '.hta', '.cmd', '.vbe')
              or
            (file.is_exec or file.is_dll)
          )
            and
        file.name imatches startup_locations
            and
            not
        ps.exe imatches
          (
            '?:\\Windows\\System32\\wuauclt.exe',
            '?:\\Windows\\System32\\msiexec.exe',
            '?:\\Windows\\SysWOW64\\msiexec.exe',
            '?:\\Windows\\System32\\svchost.exe',
            '?:\\ProgramData\\Microsoft\\Windows Defender\\Platform\\*.exe'
          )
      action: >
        {{
            emit . "Unusual file written or modified in Startup folder" ""
        }}
    - name: Unusual process modified the registry run key
      description: |
        Identifies an attempt by unusual Windows native processes to modify
        the run key and gain persistence on users logons or machine reboots.
      condition: >
        modify_registry
            and
        ps.exe imatches '?:\\Windows\\*'
            and
        registry.key.name imatches registry_run_keys
            and
            not
        ps.exe imatches
          (
            '?:\\Windows\\System32\\svchost.exe',
            '?:\\Windows\\SysWOW64\\msiexec.exe',
            '?:\\Windows\\System32\\msiexec.exe',
            '?:\\Windows\\System32\\drvinst.exe',
            '?:\\Windows\\System32\\WinSAT.exe',
            '?:\\Windows\\System32\\reg.exe',
            '?:\\Windows\\regedit.exe',
            '?:\\Windows\\SysWOW64\\reg.exe',
            '?:\\Windows\\System32\\csrss.exe',
            '?:\\Windows\\SysWOW64\\DriverStore\\*.exe',
            '?:\\Windows\\System32\\DriverStore\\*.exe',
            '?:\\Windows\\Installer\\*.exe',
            '?:\\Windows\\explorer.exe',
            '?:\\Windows\\IMECache\\*.exe',
            '?:\\Windows\\System32\\sihost.exe',
            '?:\\Windows\\SysWOW64\\prevhost.exe',
            '?:\\Windows\\System32\\conhost.exe',
            '?:\\Windows\\System32\\taskhostw.exe'
          )
      action: >
        {{
            emit . "Unusual process modified the registry run key" ""
        }}
    - name: Network connection via startup folder executable or script
      description: |
        Identifies the execution of unsigned binary or script from the
        Startup folder followed by network inbound or outbound connection.
      condition: >
        sequence
        maxspan 5m
        by ps.uuid
          |
            (
              load_untrusted_executable
                  and
              image.name imatches startup_locations
            )
                or
            (
              load_executable
                  and
              ps.name in script_interpreters
                  and
              ps.cmdline imatches startup_locations
            )
          |
          |connect_socket or accept_socket|
      action: >
        {{
            emit . "Network connection via startup folder executable or script" ""
        }}
    - name: Suspicious persistence via registry modification
      description: |
        Adversaries may abuse the registry to achieve persistence
        by modifying the keys that are unlikely modified by legitimate
        processes.
      condition: >
        modify_registry
            and
          (
            (ps.name in script_interpreters or ps.name in ('reg.exe', 'rundll32.exe', 'regsvr32.exe'))
                or
            ps.exe imatches '?:\\Users\\Public\\*'
                or
            not (pe.is_signed or pe.is_trusted)
          )
            and
        registry.key.name imatches registry_persistence_keys
      action: >
        {{
            emit . "Suspicious persistence via registry modification" ""
        }}
    - name: Suspicious Startup shell folder modification
      description: |
        Detects when adversaries attempt to modify the default Startup
        folder path to to circumvent runtime rules that hunt for file
        creations in the default Startup folder.
      condition: >
        modify_registry
            and
        registry.key.name imatches startup_shell_folder_registry_keys
            and
            not
          (
            registry.value imatches startup_locations
                or
            registry.value imatches ('%ProgramData%\\Microsoft\\Windows\\Start Menu\\Programs\\Startup')
          )
      action: >
        {{
            emit . "Suspicious Startup shell folder modification" ""
        }}
    - name: Script interpreter host or untrusted process persistence
      description: |
        Identifies the script interpreter or untrusted process writing
        to commonly abused run keys or the Startup folder locations.
      condition: >
        (modify_registry or create_file)
            and
        (
            ps.name in script_interpreters
                or
            ps.parent.name in script_interpreters
                or
            not pe.is_trusted
        )
            and
        (
            registry.key.name imatches registry_run_keys
                or
            file.name imatches startup_locations
        )
      action: >
        {{
            emit . "Script interpreter host or untrusted process persistence" ""
        }}
